# packaging/mydombox/DEBIAN/prerm
#!/bin/sh
# prerm script for mydombox

set -e

case "$1" in
    purge|remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)

	systemctl stop mydombox
    systemctl disable mydombox.service
    
	update-rc.d mydombox remove	
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

exit 0